# debian/ganglia-nagios-bridge.cron.d
* * * * * nagios /usr/bin/ganglia-nagios-bridge
